# deploy/finance-agents-phoenix.service
[Unit]
Description=Phoenix AI Observability Server (Finance Agents)
After=network.target

[Service]
Type=simple
User=financeagents
Group=financeagents
WorkingDirectory=/opt/finance-agents
Environment=PHOENIX_PORT=6006
Environment=PHOENIX_WORKING_DIR=/opt/finance-agents/phoenix_data
ExecStart=/opt/finance-agents/venv/bin/python -m phoenix.server.main serve
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/finance-agents/logs/phoenix.log
StandardError=append:/opt/finance-agents/logs/phoenix.log

[Install]
WantedBy=multi-user.target
