# deploy/finance-agents-celery.service
# =============================================================================
# Finance Agents Platform — Celery Worker systemd Service
# =============================================================================
# Install:
#   sudo cp /opt/finance-agents/deploy/finance-agents-celery.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable finance-agents-celery
#   sudo systemctl start finance-agents-celery
# =============================================================================

[Unit]
Description=Finance Agents Celery Worker
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=forking
User=financeagents
Group=financeagents
WorkingDirectory=/opt/finance-agents
EnvironmentFile=/opt/finance-agents/.env.production

ExecStart=/opt/finance-agents/venv/bin/celery \
    -A config \
    multi start worker \
    --pidfile=/opt/finance-agents/run/celery-%%n.pid \
    --logfile=/opt/finance-agents/logs/celery-%%n.log \
    --loglevel=info \
    --queues=default,agents,reconciliation,extraction,scheduled \
    --concurrency=4 \
    --max-tasks-per-child=200 \
    --time-limit=600 \
    --soft-time-limit=540 \
    -Ofair

ExecStop=/opt/finance-agents/venv/bin/celery \
    -A config \
    multi stopwait worker \
    --pidfile=/opt/finance-agents/run/celery-%%n.pid

ExecReload=/bin/kill -s HUP $MAINPID

Restart=always
RestartSec=10

# Security hardening
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
